# ansible/grow-rootfs.service
[Unit]
Description=Resize rootfs on boot

[Service]
Type=oneshot
ExecStart=/usr/bin/grow-rootfs

[Install]
WantedBy=local-fs.target
